# deployment/supervisor.conf
[program:trading_bot]
command=/home/trader/trading_bot/.venv/bin/python /home/trader/trading_bot/main.py
directory=/home/trader/trading_bot
user=trader
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/home/trader/trading_bot/logs/main.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=PYTHONPATH="/home/trader/trading_bot"
priority=100

[program:telegram_bot]
command=/home/trader/trading_bot/.venv/bin/python /home/trader/trading_bot/deployment/telegram_bot.py
directory=/home/trader/trading_bot
user=trader
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/home/trader/trading_bot/logs/telegram_bot.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PYTHONPATH="/home/trader/trading_bot"
priority=200

[group:trading_system]
programs=trading_bot,telegram_bot
priority=999
